# roles/galaxy/traefik_certs_dumper/templates/systemd/devture-traefik-certs-dumper.service.j2
[Unit]
Description=Traefik certs dumper ({{ devture_traefik_certs_dumper_identifier }})
{% for service in devture_traefik_certs_dumper_systemd_required_services_list %}
Requires={{ service }}
After={{ service }}
{% endfor %}
DefaultDependencies=no

[Service]
Type=simple
Environment="HOME={{ devture_systemd_docker_base_systemd_unit_home_path }}"

ExecStartPre=-{{ devture_systemd_docker_base_host_command_sh }} -c '{{ devture_systemd_docker_base_host_command_docker }} stop --time={{ devture_traefik_certs_dumper_container_stop_grace_time_seconds }} {{ devture_traefik_certs_dumper_identifier }} 2>/dev/null || true'
ExecStartPre=-{{ devture_systemd_docker_base_host_command_sh }} -c '{{ devture_systemd_docker_base_host_command_docker }} rm {{ devture_traefik_certs_dumper_identifier }} 2>/dev/null || true'

# This container runs as root, because:
# - Traefik runs as root and files generated by it (acme.json) are owned by root
# - it needs to be privileged enough to chown the output files (/out) as requested
#   in `devture_traefik_certs_dumper_dumped_certificates_dir_owner` and `devture_traefik_certs_dumper_dumped_certificates_dir_group`
ExecStart={{ devture_systemd_docker_base_host_command_docker }} run \
			--rm \
			--name={{ devture_traefik_certs_dumper_identifier }} \
			--log-driver=none \
			--mount type=bind,src={{ devture_traefik_certs_dumper_bin_dir_path }},dst=/certs-dumper-bin,ro \
			--mount type=bind,src={{ devture_traefik_certs_dumper_ssl_dir_path }},dst=/in,ro \
			--mount type=bind,src={{ devture_traefik_certs_dumper_dumped_certificates_dir_path }},dst=/out \
			--tmpfs=/intermediate:rw,noexec,nosuid,size=100m \
			--tmpfs=/staging:rw,noexec,nosuid,size=100m \
			--entrypoint=/bin/sh \
			{{ devture_traefik_certs_dumper_container_image }} \
			/certs-dumper-bin/entrypoint.sh

ExecStop=-{{ devture_systemd_docker_base_host_command_sh }} -c '{{ devture_systemd_docker_base_host_command_docker }} stop --time={{ devture_traefik_certs_dumper_container_stop_grace_time_seconds }} {{ devture_traefik_certs_dumper_identifier }} 2>/dev/null || true'
ExecStop=-{{ devture_systemd_docker_base_host_command_sh }} -c '{{ devture_systemd_docker_base_host_command_docker }} rm {{ devture_traefik_certs_dumper_identifier }} 2>/dev/null || true'

Restart=always
RestartSec=30
SyslogIdentifier={{ devture_traefik_certs_dumper_identifier }}

[Install]
WantedBy=multi-user.target
